# updateAll.cron
0 0 * * * python /usr/src/top15/app/updateAll.py
