# updateAll.cron
#!/bin/bash
#PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
0 0 * * * python /usr/src/top15/app/updateAll.py
